# pi-nas/crontab
30 3 * * * /opt/backup_scripts/rsync_backup.sh
